# systemd/micmute-status-monitor.service
[Unit]
Description=Mic mute Blink(1) status monitor
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/Downloads/MicMuteStatusMonitor
ExecStart=%h/Downloads/MicMuteStatusMonitor/.venv/bin/micmute-status-monitor -o blink1
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target